# config/systemd/demo.service
[Unit]
Description=Demo

[Service]
Type=idle
ExecStart=/usr/bin/demo

[Install]
WantedBy=multi-user.target